# product/wgmanager_0.0.1-1_amd64/etc/systemd/system/wgmanager.service
[Unit]
Description=Wireguard Manager

[Service]
ExecStart=/etc/wgmanager/wgmanager
WorkingDirectory=/etc/wgmanager

[Install]
WantedBy=multi-user.target